# cookbooks/webservice/templates/default/webservice.sh.erb
### BEGIN INIT INFO
# Provides:          webservice
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

INSTALL_PATH=<%= node["webservice"]["install_dir"] %>
WS_URL="<%= node["webservice"]["webservice_url_full"] %>"
PID_FILE=$INSTALL_PATH/wsdaemon.pid
if [[ -f "$PID_FILE" ]]; then
	kill -6 $(cat $PID_FILE)
	rm -f $PID_FILE
fi
export config="host=localhost,port=8649,wireformat31x=true,process=WebService,mode=unicast,config=$INSTALL_PATH/jmxetric.xml"

java -Dws.url=$WS_URL -Djava.util.logging.config.file=$INSTALL_PATH/logging.properties -javaagent:$INSTALL_PATH/jmxetric-1.0.4.jar=$config -jar $INSTALL_PATH/webservice.jar  > $INSTALL_PATH/webservice.log 2>&1 &

pid=$!
echo $pid > $PID_FILE