# cookbooks/linux_metrics/templates/default/linux_metrics.upstart.conf.erb
#!upstart
description "Linux metrics servver"

env USER=<%= node[:linux_metrics][:user] %>
start on filesystem or runlevel [2345]
stop on shutdown


script

	exec sudo -u $USER sh -c "<%= node[:linux_metrics][:full_path] %>"

end script

